# wdm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/wdm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/wdm
-rwxr-xr-x	root/root	usr/bin/wdmLogin
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/wdm/
-rwxr-xr-x	root/root	usr/etc/wdm/GiveConsole
-rwxr-xr-x	root/root	usr/etc/wdm/TakeConsole
-rw-r--r--	root/root	usr/etc/wdm/Xaccess
-rwxr-xr-x	root/root	usr/etc/wdm/Xclients
-rw-r--r--	root/root	usr/etc/wdm/Xresources
-rw-r--r--	root/root	usr/etc/wdm/Xservers
-rwxr-xr-x	root/root	usr/etc/wdm/Xsession
-rwxr-xr-x	root/root	usr/etc/wdm/Xsetup_0
drwx------	root/root	usr/etc/wdm/authdir/
-rw-------	root/root	usr/etc/wdm/wdm-config
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/wdm.1.gz
-rw-r--r--	root/root	usr/man/man1/wdmLogin.1.gz
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/run/
drwx------	root/root	var/run/wdm/
