# qmail/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/qmaild
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
lrwxrwxrwx	root/root	usr/lib/sendmail -> /var/qmail/bin/sendmail
drwxr-xr-x	root/root	usr/sbin/
lrwxrwxrwx	root/root	usr/sbin/sendmail -> /var/qmail/bin/sendmail
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/qmail/
drwxr-xr-x	root/qmail	var/qmail/
drwxr-sr-x	alias/qmail	var/qmail/alias/
drwxr-xr-x	root/qmail	var/qmail/bin/
-rwxr-xr-x	root/qmail	var/qmail/bin/bouncesaying
-rwxr-xr-x	root/qmail	var/qmail/bin/condredirect
-rwxr-xr-x	root/qmail	var/qmail/bin/datemail
-rwxr-xr-x	root/qmail	var/qmail/bin/elq
-rwxr-xr-x	root/qmail	var/qmail/bin/except
-rwxr-xr-x	root/qmail	var/qmail/bin/forward
-rwxr-xr-x	root/qmail	var/qmail/bin/maildir2mbox
-rwxr-xr-x	root/qmail	var/qmail/bin/maildirmake
-rwxr-xr-x	root/qmail	var/qmail/bin/maildirwatch
-rwxr-xr-x	root/qmail	var/qmail/bin/mailsubj
-rwxr-xr-x	root/qmail	var/qmail/bin/pinq
-rwxr-xr-x	root/qmail	var/qmail/bin/predate
-rwxr-xr-x	root/qmail	var/qmail/bin/preline
-rwxr-xr-x	root/qmail	var/qmail/bin/qail
-rwxr-xr-x	root/qmail	var/qmail/bin/qbiff
-rwx--x--x	root/qmail	var/qmail/bin/qmail-clean
-rwx--x--x	root/qmail	var/qmail/bin/qmail-getpw
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-inject
-rwx--x--x	root/qmail	var/qmail/bin/qmail-local
-rwx------	root/qmail	var/qmail/bin/qmail-lspawn
-rwx------	root/qmail	var/qmail/bin/qmail-newmrh
-rwx------	root/qmail	var/qmail/bin/qmail-newu
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-pop3d
-rwx--x--x	root/qmail	var/qmail/bin/qmail-popup
-rwx--x--x	root/qmail	var/qmail/bin/qmail-pw2u
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-qmqpc
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-qmqpd
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-qmtpd
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-qread
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-qstat
-rws--x--x	qmailq/qmail	var/qmail/bin/qmail-queue
-rwx--x--x	root/qmail	var/qmail/bin/qmail-remote
-rwx--x--x	root/qmail	var/qmail/bin/qmail-rspawn
-rwx--x--x	root/qmail	var/qmail/bin/qmail-send
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-showctl
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-smtpd
-rwx------	root/qmail	var/qmail/bin/qmail-start
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-tcpok
-rwxr-xr-x	root/qmail	var/qmail/bin/qmail-tcpto
-rwxr-xr-x	root/qmail	var/qmail/bin/qreceipt
-rwxr-xr-x	root/qmail	var/qmail/bin/qsmhook
-rwxr-xr-x	root/qmail	var/qmail/bin/sendmail
-rwx--x--x	root/qmail	var/qmail/bin/splogger
-rwxr-xr-x	root/qmail	var/qmail/bin/tcp-env
drwxr-xr-x	root/qmail	var/qmail/boot/
-rwxr-xr-x	root/qmail	var/qmail/boot/binm1
-rwxr-xr-x	root/qmail	var/qmail/boot/binm1+df
-rwxr-xr-x	root/qmail	var/qmail/boot/binm2
-rwxr-xr-x	root/qmail	var/qmail/boot/binm2+df
-rwxr-xr-x	root/qmail	var/qmail/boot/binm3
-rwxr-xr-x	root/qmail	var/qmail/boot/binm3+df
-rwxr-xr-x	root/qmail	var/qmail/boot/home
-rwxr-xr-x	root/qmail	var/qmail/boot/home+df
-rwxr-xr-x	root/qmail	var/qmail/boot/proc
-rwxr-xr-x	root/qmail	var/qmail/boot/proc+df
drwxr-xr-x	root/qmail	var/qmail/control/
-rw-r--r--	root/root	var/qmail/control/defaultdelivery
-rw-r--r--	root/root	var/qmail/control/defaultdomain
-rw-r--r--	root/root	var/qmail/control/locals
-rw-r--r--	root/root	var/qmail/control/me
-rw-r--r--	root/root	var/qmail/control/plusdomain
-rw-r--r--	root/root	var/qmail/control/rcpthosts
drwxr-xr-x	root/qmail	var/qmail/man/
drwxr-xr-x	root/qmail	var/qmail/man/cat1/
-rw-r--r--	root/qmail	var/qmail/man/cat1/bouncesaying.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/condredirect.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/except.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/forward.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/maildir2mbox.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/maildirmake.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/maildirwatch.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/mailsubj.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/preline.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/qbiff.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/qreceipt.0
-rw-r--r--	root/qmail	var/qmail/man/cat1/tcp-env.0
drwxr-xr-x	root/qmail	var/qmail/man/cat5/
-rw-r--r--	root/qmail	var/qmail/man/cat5/addresses.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/dot-qmail.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/envelopes.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/maildir.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/mbox.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/qmail-control.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/qmail-header.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/qmail-log.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/qmail-users.0
-rw-r--r--	root/qmail	var/qmail/man/cat5/tcp-environ.0
drwxr-xr-x	root/qmail	var/qmail/man/cat7/
-rw-r--r--	root/qmail	var/qmail/man/cat7/forgeries.0
-rw-r--r--	root/qmail	var/qmail/man/cat7/qmail-limits.0
-rw-r--r--	root/qmail	var/qmail/man/cat7/qmail.0
drwxr-xr-x	root/qmail	var/qmail/man/cat8/
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-clean.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-command.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-getpw.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-inject.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-local.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-lspawn.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-newmrh.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-newu.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-pop3d.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-popup.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-pw2u.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-qmqpc.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-qmqpd.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-qmtpd.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-qread.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-qstat.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-queue.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-remote.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-rspawn.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-send.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-showctl.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-smtpd.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-start.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-tcpok.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/qmail-tcpto.0
-rw-r--r--	root/qmail	var/qmail/man/cat8/splogger.0
drwxr-xr-x	root/qmail	var/qmail/man/man1/
-rw-r--r--	root/qmail	var/qmail/man/man1/bouncesaying.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/condredirect.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/except.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/forward.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/maildir2mbox.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/maildirmake.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/maildirwatch.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/mailsubj.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/preline.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/qbiff.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/qreceipt.1.gz
-rw-r--r--	root/qmail	var/qmail/man/man1/tcp-env.1.gz
drwxr-xr-x	root/qmail	var/qmail/man/man5/
-rw-r--r--	root/qmail	var/qmail/man/man5/addresses.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/dot-qmail.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/envelopes.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/maildir.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/mbox.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/qmail-control.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/qmail-header.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/qmail-log.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/qmail-users.5.gz
-rw-r--r--	root/qmail	var/qmail/man/man5/tcp-environ.5.gz
drwxr-xr-x	root/qmail	var/qmail/man/man7/
-rw-r--r--	root/qmail	var/qmail/man/man7/forgeries.7.gz
-rw-r--r--	root/qmail	var/qmail/man/man7/qmail-limits.7.gz
-rw-r--r--	root/qmail	var/qmail/man/man7/qmail.7.gz
drwxr-xr-x	root/qmail	var/qmail/man/man8/
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-clean.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-command.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-getpw.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-inject.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-local.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-lspawn.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-newmrh.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-newu.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-pop3d.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-popup.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-pw2u.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-qmqpc.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-qmqpd.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-qmtpd.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-qread.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-qstat.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-queue.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-remote.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-rspawn.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-send.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-showctl.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-smtpd.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-start.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-tcpok.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/qmail-tcpto.8.gz
-rw-r--r--	root/qmail	var/qmail/man/man8/splogger.8.gz
drwxr-x---	qmailq/qmail	var/qmail/queue/
drwx------	qmails/qmail	var/qmail/queue/bounce/
drwx------	qmails/qmail	var/qmail/queue/info/
drwx------	qmails/qmail	var/qmail/queue/info/0/
drwx------	qmails/qmail	var/qmail/queue/info/1/
drwx------	qmails/qmail	var/qmail/queue/info/10/
drwx------	qmails/qmail	var/qmail/queue/info/11/
drwx------	qmails/qmail	var/qmail/queue/info/12/
drwx------	qmails/qmail	var/qmail/queue/info/13/
drwx------	qmails/qmail	var/qmail/queue/info/14/
drwx------	qmails/qmail	var/qmail/queue/info/15/
drwx------	qmails/qmail	var/qmail/queue/info/16/
drwx------	qmails/qmail	var/qmail/queue/info/17/
drwx------	qmails/qmail	var/qmail/queue/info/18/
drwx------	qmails/qmail	var/qmail/queue/info/19/
drwx------	qmails/qmail	var/qmail/queue/info/2/
drwx------	qmails/qmail	var/qmail/queue/info/20/
drwx------	qmails/qmail	var/qmail/queue/info/21/
drwx------	qmails/qmail	var/qmail/queue/info/22/
drwx------	qmails/qmail	var/qmail/queue/info/3/
drwx------	qmails/qmail	var/qmail/queue/info/4/
drwx------	qmails/qmail	var/qmail/queue/info/5/
drwx------	qmails/qmail	var/qmail/queue/info/6/
drwx------	qmails/qmail	var/qmail/queue/info/7/
drwx------	qmails/qmail	var/qmail/queue/info/8/
drwx------	qmails/qmail	var/qmail/queue/info/9/
drwx------	qmailq/qmail	var/qmail/queue/intd/
drwx------	qmails/qmail	var/qmail/queue/local/
drwx------	qmails/qmail	var/qmail/queue/local/0/
drwx------	qmails/qmail	var/qmail/queue/local/1/
drwx------	qmails/qmail	var/qmail/queue/local/10/
drwx------	qmails/qmail	var/qmail/queue/local/11/
drwx------	qmails/qmail	var/qmail/queue/local/12/
drwx------	qmails/qmail	var/qmail/queue/local/13/
drwx------	qmails/qmail	var/qmail/queue/local/14/
drwx------	qmails/qmail	var/qmail/queue/local/15/
drwx------	qmails/qmail	var/qmail/queue/local/16/
drwx------	qmails/qmail	var/qmail/queue/local/17/
drwx------	qmails/qmail	var/qmail/queue/local/18/
drwx------	qmails/qmail	var/qmail/queue/local/19/
drwx------	qmails/qmail	var/qmail/queue/local/2/
drwx------	qmails/qmail	var/qmail/queue/local/20/
drwx------	qmails/qmail	var/qmail/queue/local/21/
drwx------	qmails/qmail	var/qmail/queue/local/22/
drwx------	qmails/qmail	var/qmail/queue/local/3/
drwx------	qmails/qmail	var/qmail/queue/local/4/
drwx------	qmails/qmail	var/qmail/queue/local/5/
drwx------	qmails/qmail	var/qmail/queue/local/6/
drwx------	qmails/qmail	var/qmail/queue/local/7/
drwx------	qmails/qmail	var/qmail/queue/local/8/
drwx------	qmails/qmail	var/qmail/queue/local/9/
drwxr-x---	qmailq/qmail	var/qmail/queue/lock/
-rw-------	qmails/qmail	var/qmail/queue/lock/sendmutex (EMPTY)
-rw-r--r--	qmailr/qmail	var/qmail/queue/lock/tcpto
prw--w--w-	qmails/qmail	var/qmail/queue/lock/trigger
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/0/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/1/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/10/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/11/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/12/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/13/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/14/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/15/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/16/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/17/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/18/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/19/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/2/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/20/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/21/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/22/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/3/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/4/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/5/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/6/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/7/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/8/
drwxr-x---	qmailq/qmail	var/qmail/queue/mess/9/
drwx------	qmailq/qmail	var/qmail/queue/pid/
drwx------	qmails/qmail	var/qmail/queue/remote/
drwx------	qmails/qmail	var/qmail/queue/remote/0/
drwx------	qmails/qmail	var/qmail/queue/remote/1/
drwx------	qmails/qmail	var/qmail/queue/remote/10/
drwx------	qmails/qmail	var/qmail/queue/remote/11/
drwx------	qmails/qmail	var/qmail/queue/remote/12/
drwx------	qmails/qmail	var/qmail/queue/remote/13/
drwx------	qmails/qmail	var/qmail/queue/remote/14/
drwx------	qmails/qmail	var/qmail/queue/remote/15/
drwx------	qmails/qmail	var/qmail/queue/remote/16/
drwx------	qmails/qmail	var/qmail/queue/remote/17/
drwx------	qmails/qmail	var/qmail/queue/remote/18/
drwx------	qmails/qmail	var/qmail/queue/remote/19/
drwx------	qmails/qmail	var/qmail/queue/remote/2/
drwx------	qmails/qmail	var/qmail/queue/remote/20/
drwx------	qmails/qmail	var/qmail/queue/remote/21/
drwx------	qmails/qmail	var/qmail/queue/remote/22/
drwx------	qmails/qmail	var/qmail/queue/remote/3/
drwx------	qmails/qmail	var/qmail/queue/remote/4/
drwx------	qmails/qmail	var/qmail/queue/remote/5/
drwx------	qmails/qmail	var/qmail/queue/remote/6/
drwx------	qmails/qmail	var/qmail/queue/remote/7/
drwx------	qmails/qmail	var/qmail/queue/remote/8/
drwx------	qmails/qmail	var/qmail/queue/remote/9/
drwxr-x---	qmailq/qmail	var/qmail/queue/todo/
-rwxr-xr-x	root/root	var/qmail/rc
drwxr-xr-x	root/qmail	var/qmail/users/
